# nous.service
[Unit]
Description=Nous - Cognitive Architecture
After=network.target ollama.service
Wants=ollama.service

[Service]
Type=simple
User=nous
Group=nous
WorkingDirectory=/opt/nous
ExecStart=/opt/nous/nous --serve --port 3333 --trust --allow-shell
Restart=always
RestartSec=5
Environment=HOME=/opt/nous

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/nous

[Install]
WantedBy=multi-user.target
